# config/config-guard-etc.service
[Unit]
Description=Config Guard for /etc (audit)
After=local-fs.target

[Service]
Type=simple
ExecStart=/home/osso/.local/bin/config-guard audit --path /etc --config /etc/config-guard/etc.toml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
